# tools/systemd/dhtnode.service.in
[Unit]
Description=OpenDHT standalone node
After=network.target
 
[Service]
EnvironmentFile=@sysconfdir@/dhtnode.conf
ExecStart=@bindir@/dhtnode -d $DHT_ARGS
KillMode=process
Restart=on-failure
Type=forking
DynamicUser=yes

[Install]
WantedBy=multi-user.target
